# data/esmo.cron
# Add here any periodic tasks to be performed

* * * * * export SSPROOT="/var/www/authBridge";/usr/bin/php $SSPROOT/modules/esmo/lib/CMHandler.php >> /dev/stderr
* * * * * export SSPROOT="/var/www/queryBridge";/usr/bin/php $SSPROOT/modules/esmo/lib/CMHandler.php >> /dev/stderr

#